# etc/cron.d/openmediavault-abcde
* * * * * root /bin/touch /dev/s[rg][0-9]*
