# etc/udev/rules.d/99-openmediavault-abcde.rules
SUBSYSTEMS=="scsi", KERNEL=="s[rg][0-9]*", RUN+="/usr/sbin/omv-abcde"